# client/com.tokenfold.push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tokenfold.push</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/usr/local/bin/tokenfold-push.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TOKENFOLD_URL</key>
        <string>https://your-server.example.com</string>
        <key>TOKENFOLD_API_KEY</key>
        <string>CHANGEME</string>
    </dict>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/tmp/tokenfold-push.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tokenfold-push.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
